# systemd/jeeves-backup.service
[Unit]
Description=Jeeves Database Backup
Documentation=https://github.com/Jeeves-Cluster-Organization/jeeves-core
After=jeeves.service

[Service]
Type=oneshot
User=%i
Group=%i

# Working directory
WorkingDirectory=/opt/jeeves

# Environment
EnvironmentFile=-/opt/jeeves/.env
Environment="BACKUP_DIR=%h/.local/share/jeeves/backups"
Environment="DATABASE_PATH=/opt/jeeves/jeeves.db"
Environment="BACKUP_RETENTION_DAYS=7"

# Execute backup script
ExecStart=/opt/jeeves/scripts/backup_database.sh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jeeves-backup

# Security
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=%h/.local/share/jeeves
ReadOnlyPaths=/opt/jeeves

[Install]
WantedBy=multi-user.target
